# neptune-screen.service
[Unit]
Description=Neptune 3 DGUS Screen Daemon for Klipper
After=klipper.service
Wants=klipper.service

[Service]
Type=simple
User=NEPTUNE_USER
ExecStart=/usr/bin/python3 NEPTUNE_INSTALL_DIR/neptune_screen.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
